# DOP_automation_2019/roles/result/files/result.service
[Unit]
Description=Result Application

[Service]
WorkingDirectory=/root/app/result
ExecStart=node --require dotenv/config server.js
Restart=always

[Install]
WantedBy=multi-user.target
